# etc/supervisor/conf.d/envoy_config.ini
[program:envoy_config]
command=/usr/local/bin/configure.py /etc/envoy/
startsecs=0
priority=900
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
